# prototype/systemd/nextcloud-webhook.service
[Unit]
Description=Nextcloud webhook listener (server push sync)
After=default.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/nextcloud-server/scripts/nextcloud-webhook-listener.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
