# prototype/systemd/nextcloud-sync.service
[Unit]
Description=Nextcloud CLI sync (fallback)

[Service]
Type=oneshot
ExecStart=%h/nextcloud-server/scripts/nextcloud-sync.sh
